# supervisord.conf
[supervisord]
nodaemon=true

[program:haproxy]
command=haproxy -f /usr/local/etc/haproxy/haproxy.cfg
